# cfp-demo-car/systemd/cfp-obd.service
[Unit]
Description=CFP OBD2 CAN Bus Emulator
After=network.target
Wants=network.target

[Service]
ExecStartPre=/sbin/ip link set can1 up type can bitrate 500000
ExecStart=/usr/bin/python3 /home/s10skeleton/CFPDemo/cfp-demo-car/obd_emulator.py
Environment=CFP_SIMULATE=0
Restart=always
RestartSec=3
User=s10skeleton
WorkingDirectory=/home/s10skeleton/CFPDemo/cfp-demo-car

[Install]
WantedBy=multi-user.target
